# systemd/my-owntracks-digest.service
[Unit]
Description=Telegram Control OwnTracks daily digest

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/projects/telegram_control
Environment=TZ=Asia/Kolkata
ExecStart=%h/.local/bin/uv --cache-dir .uv-cache run python -m owntracks.digest --date today --send
