# systemd/my-owntracks-digest.service
[Unit]
Description=Telegram Control OwnTracks daily digest

[Service]
Type=oneshot
WorkingDirectory=/path/to/telectl
Environment=TZ=Asia/Kolkata
ExecStart=%h/.local/bin/uv --cache-dir .uv-cache run python -m owntracks.digest --date today --send
